# Scripts/99-usb.rules
ACTION=="add", SUBSYSTEM=="block" ATTRS{idVendor}=="****", ATTRS{idProduct}=="****", RUN+="/usr/local/bin/firmware-update.sh"
